# crontab
* * * * * root cd /app && python main.py >> /var/log/cron.log 2>&1
# @hourly root cd /app && python main.py >> /var/log/cron.log 2>&1
